# packer/rootfs/etc/systemd/system/udptunnel.socket
[Unit]
Description=udptunnel socket

[Socket]
ListenStream=443
BindIPv6Only=both

[Install]
WantedBy=sockets.target
